# systemd/sennet-portal.service
[Unit]
Description=SenNet Portal (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sennet
WorkingDirectory=/opt/sennet-portal/portal
ExecStart=/usr/bin/env HOSTNAME=0.0.0.0 PORT=3000 node /opt/sennet-portal/portal/.next/standalone/server.js
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
