# systemd/sennet-scheduler-worker.service
[Unit]
Description=SenNet Scheduler Worker (FastAPI-only oneshot)
After=network-online.target sennet-agent-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=debian
WorkingDirectory=/opt/sennet-agent/repo
Environment="PATH=/opt/sennet-agent/venv/bin"
ExecStart=/opt/sennet-agent/venv/bin/python /opt/sennet-agent/repo/app/run_report_oneshot.py
