# systemd/sennet-agent.service
[Unit]
Description=DEPRECATED SenNet Streamlit UI (disabled in FastAPI-only architecture)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=debian
WorkingDirectory=/opt/sennet-agent/repo/app
Environment="PATH=/opt/sennet-agent/venv/bin"
ExecStart=/bin/false

[Install]
WantedBy=multi-user.target
